# dot_config/systemd/user/yubikey-agent.service.d/override.conf
[Service]
Environment=SSH_AUTH_SOCK=%t/yubikey-agent/yubikey-agent.sock
Environment=DISPLAY=:0
Environment=GPG_TTY=/dev/pts/1
Environment="XAUTHORITY=/home/max/.Xauthority"
